# ctrl.conf
[program:ctrl_script]

command=python /home/ats/project/ctrl.py

autostart=true

autorestart=true

stderr_logfile=/var/log/ctrl.err.log

stdout_logfile=/var/log/ctrl.out.log
